# flask_app.service
[Unit]
Description=Gunicorn instance to serve Flask app
After=network.target

[Service]
User=your_user
Group=www-data
WorkingDirectory=/path/to/your/app
ExecStart=/usr/local/bin/gunicorn -c /path/to/your/app/gunicorn_conf.py app:app

[Install]
WantedBy=multi-user.target
